# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:portal]
command=gunicorn -b 0.0.0.0:8000 run
directory=/root/portal
redirect_stderr=True

