# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:portal]
command=python /root/portal/run.py --port=8080

