# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=*:9001

[program:portal]
command=bash -c 'cd /root/portal && gunicorn -b 0.0.0.0:8000 run'

